# ssdpd.service.in
[Unit]
Description=Simple Service Discovery Protocol daemon (UPnP)
Documentation=man:ssdpd
Documentation=https://github.com/troglobit/ssdpd
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=@SBINDIR@/ssdpd -sn

[Install]
WantedBy=multi-user.target
